# systemd/nbd-vram-battery-watch.service
[Unit]
Description=nbd-vram battery threshold watch
After=vram-swap-nbd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nbd-vram-power-check.sh
