# systemd/nbd-vram-power-check.service
[Unit]
Description=nbd-vram power state check
After=vram-swap-nbd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nbd-vram-power-check.sh
