# systemd/vram-swap-nbd.service
[Unit]
Description=VRAM Swap via NBD (CUDA memcpy backend)
After=systemd-modules-load.service nvidia-persistenced.service
Wants=nvidia-persistenced.service

[Service]
Type=notify
NotifyAccess=main

# VRAM to allocate for swap. Leave headroom for the GPU's own use.
# RTX 3070 (8GB): 7168 leaves ~1GB for display/compute.
Environment=VRAM_SETUP_SIZE_MB=7168
Environment=VRAM_SWAP_PRIORITY=1500

# Wake GPU before CUDA init to avoid GSP RPC timeout (Xid 119 on laptops)
ExecStartPre=/usr/bin/nvidia-smi -pm 1

ExecStart=/usr/local/bin/nbd-vram
ExecStartPost=/usr/local/bin/nbd-vram-connect.sh
ExecStop=/usr/local/bin/nbd-vram-disconnect.sh

Restart=on-failure
RestartSec=10
TimeoutStartSec=60
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbd-vram

[Install]
WantedBy=multi-user.target
